# userpatches/overlay/config_files/system/_etc_systemd_system_rnsd.service
[Unit]
Description=RNS Server

[Service]
Type=simple
User=__USER__
ExecStart=/home/__USER__/.local/bin/rnsd

[Install]
WantedBy=multi-user.target
